# system_files/usr/lib/systemd/system/fleet-flatpak-setup.timer
# Déclenche la synchronisation Flatpak QUOTIDIENNE (MAJ --or-update + prune des apps
# retirées des listes). Le réessai rapproché du premier boot est géré par le service
# lui-même (Restart=on-failure, voir fleet-flatpak-setup.service) — pas par ce timer.
# Persistent=true : rattrape une échéance manquée (poste éteint) au démarrage suivant.
[Unit]
Description=Synchronisation quotidienne des Flatpaks de la flotte

[Timer]
OnCalendar=daily
Persistent=true
RandomizedDelaySec=1h

[Install]
WantedBy=timers.target
