# system_files/usr/lib/systemd/system/fleet-admin-desktop-setup.service
# Dépose le script change_password.sh à la racine du home d'admin, au premier boot.
# Activé au build (build.sh). oneshot, une seule fois (stamp posé après).
[Unit]
Description=Dépose change_password.sh dans le home d'admin (1er boot)
ConditionPathExists=!/var/lib/fleet/admin-home-done
After=multi-user.target

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStartPre=/usr/bin/mkdir -p /var/lib/fleet
# Si admin existe : installe le script (exécutable, possédé par admin) à la racine du home.
ExecStart=/usr/bin/bash -c 'id admin >/dev/null 2>&1 && install -o admin -g admin -m 0755 /usr/share/fleet/change_password.sh /home/admin/change_password.sh || true'
ExecStartPost=/usr/bin/touch /var/lib/fleet/admin-home-done

[Install]
WantedBy=multi-user.target
