# system_files/usr/lib/systemd/system/fleet-flatpak-setup.service
# Installe/synchronise les Flatpaks de la flotte.
# DÉCLENCHÉ PAR fleet-flatpak-setup.timer (réessai chaque minute jusqu'au succès).
# Pas de RemainAfterExit : le service retourne à l'état "inactive" après chaque tentative,
# ce qui permet au timer (OnUnitInactiveSec) de replanifier. Le stamp de succès
# /var/lib/fleet/flatpaks-installed coupe les relances (Condition ci-dessous).
[Unit]
Description=Installation des Flatpaks de la flotte
Wants=network-online.target
After=network-online.target
ConditionPathExists=/usr/share/fleet/flatpaks.list
ConditionPathExists=!/var/lib/fleet/flatpaks-installed

[Service]
Type=oneshot
ExecStart=/usr/libexec/fleet-flatpak-sync
# Gros téléchargements possibles (LibreOffice, VS Code…) : timeout généreux.
TimeoutStartSec=1800

# Pas de [Install] : ce service n'est pas lancé directement au boot, c'est le timer
# (fleet-flatpak-setup.timer) qui le déclenche. Voir build.sh (enable du timer).
