# system_files/usr/lib/systemd/system/fleet-tpm-enroll.service
[Unit]
Description=Enrolement TPM2 du deverrouillage LUKS (premier boot)
# Ne s'execute qu'une fois : un stamp est pose au premier passage.
ConditionPathExists=!/var/lib/fleet/tpm-enrolled
After=multi-user.target
Before=graphical.target

[Service]
Type=oneshot
RemainAfterExit=yes
# Marque "tente" AVANT le lancement -> pas de re-prompt a chaque boot meme si annule.
ExecStartPre=/usr/bin/mkdir -p /var/lib/fleet
ExecStartPre=/usr/bin/touch /var/lib/fleet/tpm-enrolled
ExecStart=/usr/libexec/fleet-tpm-enroll
# Prompt de la passphrase LUKS existante sur la console, au premier demarrage.
StandardInput=tty-force
StandardOutput=tty
StandardError=tty
TTYPath=/dev/console

[Install]
WantedBy=multi-user.target
