# system_files/usr/lib/systemd/system/fleet-initramfs-localize.service
# Active la régénération LOCALE de l'initramfs au premier boot, pour embarquer le keymap
# clavier suisse (/etc/vconsole.conf → invite LUKS/PIN en suisse, pas en QWERTY US).
# On NE régénère PAS dans le build conteneur (dracut y échoue) : on délègue à rpm-ostree
# côté poste, où la régénération est fiable (validé).
#
# Effet : crée un nouveau déploiement avec l'initramfs régénéré → appliqué au PROCHAIN
# redémarrage (pas de reboot automatique pour ne pas interrompre la mise en service).
# oneshot, une seule fois (stamp). Idempotent (n'active que si pas déjà actif).
[Unit]
Description=Régénération locale de l'initramfs (clavier suisse à l'invite LUKS) au 1er boot
ConditionPathExists=!/var/lib/fleet/initramfs-localized
After=multi-user.target

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStartPre=/usr/bin/mkdir -p /var/lib/fleet
ExecStart=/usr/bin/bash -c 'rpm-ostree initramfs 2>/dev/null | grep -qiE "regeneration: enabled|enabled" || rpm-ostree initramfs --enable'
ExecStartPost=/usr/bin/touch /var/lib/fleet/initramfs-localized

[Install]
WantedBy=multi-user.target
