# system_files/usr/lib/systemd/system/fleet-flatpak-setup.service
# Installe/synchronise les Flatpaks de la flotte (listes /usr/share/fleet/flatpaks.d/).
# Activé au build (build.sh). Déclenché : au boot (multi-user.target) ET chaque jour
# par fleet-flatpak-setup.timer (postes allumés longtemps → MAJ Flatpak régulières).
# PAS de RemainAfterExit : l'unité redevient inactive après exécution, sinon le timer
# ne pourrait pas la redéclencher.
[Unit]
Description=Installation/synchronisation des Flatpaks de la flotte
Wants=network-online.target
After=network-online.target NetworkManager-wait-online.service
# Flathub a besoin du réseau ; on ne bloque pas le boot si ça échoue.
ConditionPathExists=/usr/share/fleet/flatpaks.d

[Service]
Type=oneshot
ExecStart=/usr/libexec/fleet-flatpak-sync
# Réessaie au prochain boot / prochain déclenchement du timer si le réseau manquait.
TimeoutStartSec=infinity

[Install]
WantedBy=multi-user.target
